# luci-app-fchomo/root/etc/uci-defaults/luci-app-fchomo
#!/bin/sh

if ! uci -q get fchomo.global.authentication; then
	uci add_list fchomo.global.authentication="fchomodef:$(cat /proc/sys/kernel/random/uuid)"
	uci commit fchomo
fi

exit 0
